# webvirtmgr/webvirtmgrd.conf
[program:webvirtmgr]
command=/data/webvirtmgr/manage.py run_gunicorn -c /data/webvirtmgr/conf/gunicorn.conf.py 

[program:webvirtmgr-console]
command=/data/webvirtmgr/console/webvirtmgr-console
