# gitlab/gitlab.conf
[program:gitlab]
command=/opt/gitlab/embedded/bin/runsvdir-start 
