# hexo/hexod.conf
[program:hexo]
directory=/data/hexo
command=hexo server
