# hexo/hexod.conf
[program:hexo]
directory=/data/hexo
command=hexo server
stdout_logfile=/hexo.log
stderr_logfile=/hexo_err.log
